# data/systemd/plank-reloaded.service
[Unit]
Description=Plank Reloaded
After=graphical-session.target bamfdaemon.service
BindsTo=graphical-session.target
Wants=bamfdaemon.service

[Service]
Type=simple
TimeoutStartSec=5s
ExecStartPre=/bin/sh -c 'while ! busctl --user status org.ayatana.bamf >/dev/null 2>&1; do sleep 0.5; done'
ExecStart=plank
Restart=on-failure
RestartSec=1s
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
